# server/ente-sync.timer
[Unit]
Description=Run Ente photo export hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
