# server/ente-sync.service
[Unit]
Description=Ente photo export

[Service]
Type=oneshot
EnvironmentFile=/srv/backups/server/.env
ExecStart=/srv/backups/bin/ente-sync
